# scripts/linux-device-read/bme280-json.service
# Install on the Pi (see NEW_PI_SETUP.md in repo root):
#   sudo cp ~/Temp-Humidity-Project/scripts/linux-device-read/bme280-json.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now bme280-json.service
#
# If your Linux user is not "robot", change User=, Group=, WorkingDirectory=, and
# /home/robot paths in ExecStart= before copying to /etc/systemd/system/.
#
# Do NOT add SupplementaryGroups=i2c gpio unless both groups exist (getent group i2c gpio).
# Missing "gpio" on Ubuntu causes systemd exit status 216/GROUP.

[Unit]
Description=BME280 readings as JSON (linux-device-read)
After=network.target

[Service]
Type=simple
User=robot
Group=robot
# Do not list SupplementaryGroups= unless those groups exist (getent group i2c gpio).
# User=robot already gets groups from: sudo usermod -aG i2c robot
StateDirectory=bme280
Environment=BME280_JSON_PATH=/var/lib/bme280/readings.json
WorkingDirectory=/home/robot
ExecStart=/home/robot/bme280-venv/bin/python /home/robot/Temp-Humidity-Project/scripts/linux-device-read/bme280_daemon.py --interval 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
