# scripts/linux-device-read/bme280-json.service
# Install on the Pi:
#   sudo cp ~/scripts/linux-device-read/bme280-json.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now bme280-json.service
#
# Adjust User= and paths if your account name or layout differs from ``robot``.

[Unit]
Description=BME280 readings as JSON (linux-device-read)
After=network.target

[Service]
Type=simple
User=robot
Group=robot
SupplementaryGroups=i2c gpio
StateDirectory=bme280
Environment=BME280_JSON_PATH=/var/lib/bme280/readings.json
WorkingDirectory=/home/robot
ExecStart=/home/robot/bme280-venv/bin/python /home/robot/scripts/linux-device-read/bme280_daemon.py --interval 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
